# systemd/cloudflared-openclawdbot.service
[Unit]
Description=cloudflared (OpenClawDBot public tunnel)
After=network-online.target openclawdbot.service
Wants=network-online.target
Requires=openclawdbot.service
PartOf=orchestrator.service openclawdbot.service
BindsTo=openclawdbot.service

[Service]
TimeoutStartSec=15
Type=notify
EnvironmentFile=-/home/oneclickwebsitedesignfactory/.openclaw/workspace/orchestrator/.env
ExecStart=/usr/bin/cloudflared --no-autoupdate tunnel run --token ${OPENCLAWDBOT_TUNNEL_TOKEN}
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=default.target
